# HW17/roles/client/templates/borg-backup.timer.j2
[Unit]
Description=Borg Backup

[Timer]
OnUnitActiveSec={{ backup_frequency }}

[Install]
WantedBy=timers.target
